# recipes-example/read-load-avg/files/myservice.service
[Unit]
Description=My service

[Service]
User=monitoring
ExecStart=/usr/bin/log_loadavg.sh
